# resources/gologger/logstore.service
[Unit]
Description=Modell Aachen log persister
After=network.target

[Service]
EnvironmentFile=/etc/default/logstore
Environment=RABBIT_USER=
Environment=RABBIT_PASSWORD=
Environment=POSTGRES_USER=
Environment=POSTGRES_PASSWORD=
User=www-data
ExecStart=/usr/local/bin/logstore
Restart=on-failure

[Install]
WantedBy=multi-user.target
